# config/voicebox.service
[Unit]
Description=VoiceBox Dynamic Voice Recognition System
After=graphical.target sound.service
Wants=sound.service

[Service]
Type=simple
User=varun
Group=varun
WorkingDirectory=/home/varun/voicebox
Environment="DISPLAY=:0"
Environment="PULSE_RUNTIME_PATH=/run/user/1000/pulse"
Environment="XDG_RUNTIME_DIR=/run/user/1000"
ExecStart=/home/varun/voicebox/venv/bin/python /home/varun/voicebox/src/voicebox.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Resource limits
MemoryMax=512M
CPUQuota=80%

[Install]
WantedBy=graphical.target