# config/voicebox-kiosk.service
[Unit]
Description=VoiceBox kiosk (bare X11)
After=network.target

[Service]
User=varun
Environment=HOME=/home/varun
ExecStart=/usr/bin/startx -- -nocursor
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target